# uwsgi/farm-api-supervisord.conf
[program:farm_api]
command=/root/miniconda3/envs/farm-api/bin/gunicorn -w4 -b0.0.0.0:9000 farm_api.wsgi:application
directory=/root/farm_api/farm_api
autostart=true
autorestart=true
