# uwsgi/farm-api-supervisord.ini
[program:farm_api]
command=/root/miniconda3/envs/farm-api/bin/uwsgi --ini /root/farm_api/uwsgi/uwsgi.ini
numprocs=1
stdout_logfile=/root/farm_api/uwsgi/out.log
stderr_logfile=/root/farm_api/uwsgi/err.log
autostart=true
autorestart=true
startsecs=10
priority=997
