# uwsgi/nginx-supervisord.conf
[program:nginx]
command=/usr/sbin/nginx
startsecs=0
stopwaitsecs=0
autostart=false
autorestart=false
stdout_logfile=/home/farm-api/uwsgi/nginx.log
stderr_logfile=/home/farm-api/uwsgi/nginx.err